# init/systemd/system/glint.service.d/override.conf
[Service]
Environment="GLINTSERVER_CONFIG_FILE=/etc/glint/glintserver.conf"
